# roles/hack/files/okd.service
[Unit]
Description=OpenShift

[Service]
ExecStart=/opt/origin/cluster_start

[Install]
WantedBy=multi-user.target